# sv/blkmond/run
#!/bin/sh
cd ~/g/pushpool
exec ./blkmond ~/blkmon.cfg 2>&1
