# sv/namecoind/run
#!/bin/sh
cd ~
exec ./bin/namecoind 2>&1
