# sv/pushpoold/run
#!/bin/sh
cd ~
ulimit -n 65535
exec ./sbin/pushpoold -F -E -D2 2>&1
