# sv/pushpoold/run
#!/bin/sh
cd ~
exec ./sbin/pushpoold 2>&1
