# sv/pushpoold/run
#!/bin/sh
cd ~
ulimit -n 65535
ulimit -n > /tmp/ulimit.$$.$(date +%s)
exec ./sbin/pushpoold -F -E -D1 2>&1
